# scraping_service/packaging/debian/DEBIAN/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "upgrade" ] || [ "$1" = "deconfigure" ]; then
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop scraping-service.service || true
        systemctl disable scraping-service.service || true
    fi
fi

exit 0
